# systemd/watchpot.service
[Unit]
Description=WatchPot Photo Capture and Email Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
Group=pi
ExecStart=/opt/watchpot/scripts/watchpot.sh
WorkingDirectory=/opt/watchpot
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/watchpot /var/lib/watchpot

# Restart policy
Restart=Always
RestartSec=3600

[Install]
WantedBy=multi-user.target
